# deploy/rack-kiosk.service
[Unit]
Description=TURZX Kiosk — drives the Turing 8.8" USB display (apps + Grafana dashboard + HA/MQTT)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
# Runs as root for direct libusb access to the vendor-specific USB display.
User=root
WorkingDirectory=/home/kiosk/turzx
# Optional env (MQTT broker creds, dwell, brightness) — see kiosk.env.example
EnvironmentFile=-/home/kiosk/turzx/kiosk.env
ExecStart=/home/kiosk/turzx/venv/bin/python /home/kiosk/turzx/kiosk.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
